# docker/supervisord.conf
[program:php]
command=php-fpm

[program:nginx]
command=nginx -g 'daemon off;'